# src/etc/systemd/system/certbot.service.d/10-override.conf
[Service]
Group=vyattacfg
ExecStart=/usr/libexec/vyos/op_mode/pki.py renew_certbot
